# target/linux/avanta/base-files/etc/uci-defaults/01_leds
#!/bin/sh
#
# Copyright (C) 2012-2014 OpenWrt.org
#

. /lib/functions/uci-defaults.sh
. /lib/kirkwood.sh

board=$(kirkwood_board_name)

case "$board" in
"f660")
	ucidef_set_led_default "warn" "warn" "f660:red:warn" "1"
	ucidef_set_led_default "voice1" "voice1" "f660:green:voice1" "1"
	ucidef_set_led_default "ponnetwork" "ponnetwork" "f660:green:ponnetwork" "1"
	ucidef_set_led_default "ponport" "ponport" "f660:red:ponport" "1"
	ucidef_set_led_default "voice2" "voice2" "f660:green:voice2" "1"
	ucidef_set_led_default "wpsred" "wpsred" "f660:red:wps" "1"
	ucidef_set_led_default "wpsgreen" "wpsgreen" "f660:green:wps" "1"
	ucidef_set_led_default "usb" "usb" "f660:green:usb" "1"
	;;
*)
	;;
esac

ucidef_commit_leds

exit 0
